# components/cf-release/conf/supervisord.conf
[supervisord]
nodaemon=true

[program:monit]
command=/root/control/monit_daemon.sh
